# logrotate/files/logrotate.cron
## This file is managed by puppet

# execute logrotate at midnight
0  0  *  *  *    root  /usr/local/sbin/run-logrotate.sh
